# board/tdc/cvr-mil-v2-a/rootfs_overlay/etc/profile.d/qt.sh
export FB_MULTI_BUFFER=4
export QT_EGLFS_IMX6_NO_FB_MULTI_BUFFER=0
export QT_QPA_EGLFS_FB=/dev/fb1
export QT_QPA_EVDEV_KEYBOARD_PARAMETERS=Linuxinput:/dev/input/event0:keymap=/etc/cvr/atm88pa.qmap
export QT_QPA_FONTDIR=/usr/share/fonts
export GST_EVR_OVL_MASK=13

