# board/tdc/etlmd-v2-lp/rootfs_overlay/etc/cron.hourly/histroy-cleanup
#!/bin/sh

/etc/init.d/S90history-cleanup start
