# board/tdc/hinvr-hi3520d/rootfs_overlay/etc/profile.d/qt.sh
#!/bin/sh

export QT_QWS_FONTDIR=/usr/share/fonts
export QWS_DISPLAY=LinuxFB:/dev/fb0:nographicsmodeswitch
export QWS_MOUSE_PROTO=MouseMan:/dev/input/mice

