# board/tdc/hivcap/rootfs_overlay/etc/profile.d/gpg.sh
#!/bin/sh

export GNUPGHOME=/etc/gnupg

