# board/tdc/cvr-mil-v3/rootfs_overlay/etc/profile.d/qt.sh
export FB_MULTI_BUFFER=2
export QT_EGLFS_IMX6_NO_FB_MULTI_BUFFER=0

