# board/tdc/cvr-mil-v1/rootfs_overlay/etc/profile.d/qt.sh
#!/bin/sh

export QT_QWS_FONTDIR=/usr/share/fonts
export QWS_DISPLAY=LinuxFb:/dev/fb0
